# etc/crontab
*/1 * * * *     root    cd /root/rube-files && /root/scripts/wrapper
